# FotoAlert/deploy/fotoalert-precompute.timer
[Unit]
Description=FotoAlert – Täglich 00:01 Uhr Precompute ausführen
Requires=fotoalert-precompute.service

[Timer]
# Täglich um 00:01 Uhr Serverzeit (UTC) → 02:01 Uhr MEZ / 02:01 Uhr MESZ
OnCalendar=*-*-* 00:01:00
# Leichte Zufallsverzögerung (bis 5 Min.) verhindert Server-Lastspitzen bei mehreren Timern
RandomizedDelaySec=300
# Falls Server um 05:30 ausgeschaltet war, beim nächsten Start nachholen
Persistent=true

[Install]
WantedBy=timers.target
